# ros-sensors/debian/postinst
#!/bin/bash

systemctl daemon-reload

systemctl enable ros-sensors.service
systemctl start ros-sensors.service