# ros-sensors/debian/prerm
#!/bin/bash

systemctl stop ros-sensors.service
systemctl disable ros-sensors.service

systemctl daemon-reload
