# omron-daemon/debian/postrm
#!/bin/bash

systemctl stop omron-driver.service
systemctl disable omron-driver.service

systemctl daemon-reload
