# omron-daemon/debian/postinst
#!/bin/bash

systemctl daemon-reload

systemctl enable omron-driver.service
systemctl start omron-driver.service --no-block

