# ros-daemon/debian/postinst
#!/bin/bash

systemctl daemon-reload

systemctl enable ros.service
systemctl start ros.service