# ros-daemon/debian/prerm
#!/bin/bash

systemctl stop ros.service
systemctl disable ros.service

systemctl daemon-reload
