# robot-bringup-daemon/debian/postinst
#!/bin/bash

systemctl daemon-reload

systemctl enable robot-bringup.service
systemctl start robot-bringup.service --no-block

