# robot-bringup-daemon/debian/prerm
#!/bin/bash

systemctl stop robot-bringup.service
systemctl disable robot-bringup.service

systemctl daemon-reload
