# ros-project/debian/postinst
#!/bin/bash

systemctl daemon-reload

systemctl enable ros-project.service
systemctl start ros-project.service