# ros-project/debian/prerm
#!/bin/bash

systemctl stop ros-project.service
systemctl disable ros-project.service

systemctl daemon-reload
