# panda-daemon/debian/prerm
#!/bin/bash

systemctl stop panda-driver.service
systemctl disable panda-driver.service

systemctl daemon-reload
