# panda-realsense-daemon/debian/prerm
#!/bin/bash

systemctl stop panda-realsense.service
systemctl disable panda-realsense.service

systemctl daemon-reload
