# ros-watchdog/debian/postinst
#!/bin/bash

systemctl daemon-reload

systemctl enable ros-watchdog.service
systemctl start ros-watchdog.service --no-block
