# ros-watchdog/debian/prerm
#!/bin/bash

systemctl stop ros-watchdog.service
systemctl disable ros-watchdog.service

systemctl daemon-reload
