# ros-watchdog/etc/ros-watchdog.service.in
[Unit]
Description=ROS Watchdog Daemon
Requires=ros.service
After=ros.service

[Service]
Type=notify
ExecStart=/bin/bash -c "source /etc/qcr/default_config && source $ROS_WORKSPACE && /usr/bin/ros-watchdog"
Restart=always
TimeoutStartSec=infinity

[Install]
WantedBy=ros.service
