# roscore-daemon/debian/postrm
#!/bin/bash

systemctl stop roscore.service
systemctl disable roscore.service

systemctl daemon-reload
