# roscore-daemon/debian/postinst
#!/bin/bash

systemctl daemon-reload

systemctl enable roscore.service
systemctl start roscore.service --no-block

