# ros-robot/debian/postinst
#!/bin/bash

systemctl daemon-reload

systemctl enable ros-robot.service
systemctl start ros-robot.service