# ros-robot/debian/prerm
#!/bin/bash

systemctl stop ros-robot.service
systemctl disable ros-robot.service

systemctl daemon-reload
