# qcr-services/etc/roscore.service
[Unit]
Description=ROSCore
Requires=ros.service
After=ros.service

[Service]
ExecStart=/bin/bash -c "source /etc/qcr/qcr-env.bash && source $ROS_WORKSPACE && [ $QCR_ROS_IP_PERIPHERAL == true ] && roscore || exit 0"
Restart=always
RestartSec=5

[Install]
WantedBy=ros.service
