# roscore-daemon/etc/roscore.service.in
[Unit]
Description=ROSCore
Requires=ros.service
After=ros.service

[Service]
ExecStart=/bin/bash -c "source /etc/qcr/default_config && source $ROS_WORKSPACE && roscore"
Restart=always
RestartSec=5

[Install]
WantedBy=ros.service
